# deploy/confs/supervisor.docker.conf
[program:uwsgi]
directory=/srv/manti.by/src/app
command=uwsgi --ini /etc/manti.by/uwsgi.ini
stderr_logfile = /var/log/manti.by/uwsgi_error.log
stdout_logfile = /var/log/manti.by/uwsgi_access.log
autostart=true
autorestart=true
stopsignal=INT

[program:celery-worker]
directory=/srv/manti.by/src/app
command=celery -A core worker
stderr_logfile = /var/log/manti.by/celery_worker_error.log
stdout_logfile = /var/log/manti.by/celery_worker_access.log
autostart=true
autorestart=true
stopsignal=INT

[program:celery-flower]
directory=/srv/manti.by/src/app
command=celery flower -A core --url_prefix=flower --unix_socket=/run/manti/flower.sock
stderr_logfile = /var/log/manti.by/celery_flower_error.log
stdout_logfile = /var/log/manti.by/celery_flower_access.log
autostart=true
autorestart=true
stopsignal=INT

[group:manti]
programs=uwsgi,celery-worker,celery-flower
priority=999
